# mwan3/files/etc/init.d/mwan3
#!/bin/sh /etc/rc.common
START=99

start() {
	/usr/sbin/mwan3 mwan_apply
}

stop() {
	/usr/sbin/mwan3 mwan_unload
}

restart() {
	stop
	start
}
