# mwan3/files/etc/init.d/mwan3
#!/bin/sh /etc/rc.common
START=99

start() {

	# load /etc/config/mwan3 and run wan_ifup function on each configured and enabled mwan3 interface
	/usr/sbin/mwan3 mwan_apply &

}
