# mwan3-luci/files/etc/uci-defaults/mwan3
#!/bin/sh

uci batch <<-EOF
        add ucitrack mwan3
        set ucitrack.@mwan3[-1].exec="/usr/sbin/mwan3 mwan_apply"
        commit ucitrack
EOF
