# mwan3-luci/files/etc/uci-defaults/mwan3
#!/bin/sh

uci -q uci batch <<-EOF >/dev/null
		delete ucitrack.@mwan3[-1]
		add ucitrack mwan3
		set ucitrack.@mwan3[-1].init=mwan3
		commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
